# deploy/mannenvakanties.service
# systemd-service voor Mannenvakanties
# Plaats dit bestand als:  /etc/systemd/system/mannenvakanties.service
# Pas de regels met << >> aan jouw situatie aan.

[Unit]
Description=Mannenvakanties album
After=network.target

[Service]
Type=simple
# Een eigen gebruiker (geen root). Maak die eventueel met:  sudo useradd -r -s /usr/sbin/nologin mannen
User=<<gebruiker, bv. mannen>>
WorkingDirectory=<<pad, bv. /var/www/mannenvakanties>>
# Controleer het pad naar node met:  which node
# Dependencies bijwerken vóór elke start (bv. na een 'git pull' die package.json wijzigde).
# Het '-' ervoor betekent: faalt dit, dan start de app gewoon toch door — de site gaat nooit plat.
ExecStartPre=-/usr/bin/env npm install --omit=dev --no-audit --no-fund
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=5
# De eerste install van de detectiebibliotheken (~285 MB) kan even duren; ruim de tijd geven.
TimeoutStartSec=600

Environment=NODE_ENV=production
Environment=PORT=3000
# npm-cache in de (afgeschermde) /tmp, zodat install ook lukt voor een systeemgebruiker zonder home-map.
Environment=npm_config_cache=/tmp/npm-cache
# Zet hier een lange, willekeurige geheime string (bv. uit:  openssl rand -hex 32 )
Environment=SESSION_SECRET=<<lange-willekeurige-string>>

# Wat extra afscherming
NoNewPrivileges=true
ProtectSystem=full
PrivateTmp=true

[Install]
WantedBy=multi-user.target
